# deploy/tessitura-next.service
[Unit]
Description=Tessitura Next.js Application
After=network.target postgresql.service

[Service]
Type=simple
User=tessitura
Group=tessitura
WorkingDirectory=/home/tessitura/app
ExecStart=/usr/bin/node /home/tessitura/app/node_modules/.bin/next start -p 3000
Restart=always
RestartSec=5
Environment=NODE_ENV=production
EnvironmentFile=/home/tessitura/.env

[Install]
WantedBy=multi-user.target
